# deploy/athena.service
[Unit]
Description=Athena Weekly Planner
After=network.target

[Service]
Type=notify
User=pi
Group=pi
WorkingDirectory=/opt/athena/app
Environment="PATH=/opt/athena/venv/bin"
EnvironmentFile=/opt/athena/app/.env
ExecStart=/opt/athena/venv/bin/gunicorn -c /opt/athena/gunicorn.conf.py app:app
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/athena
ReadWritePaths=/opt/duobrain/vault

[Install]
WantedBy=multi-user.target
